# packages/agentopera/agentopera-0.0.5.tar.gz/agentopera-0.0.5/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
loglevel=info
pidfile=/tmp/supervisord.pid

[program:grpc_runtime]
command=python app/host_runtime.py  # No need for /app/app/
directory=/app
autostart=true
autorestart=true
startsecs=5
stderr_logfile=/var/log/supervisor/grpc_runtime.err.log
stdout_logfile=/var/log/supervisor/grpc_runtime.out.log
stopasgroup=true
killasgroup=true

[program:fastapi_server]
command=uvicorn app.main:app --host 0.0.0.0 --port 8000 --reload  # No app.main since no 'app' package
directory=/app
autostart=true
autorestart=true
startsecs=5
stderr_logfile=/var/log/supervisor/fastapi_server.err.log
stdout_logfile=/var/log/supervisor/fastapi_server.out.log
stopasgroup=true
killasgroup=true