# packages/rox-septentrio/rox_septentrio-0.5.2.tar.gz/rox_septentrio-0.5.2/install/ntrip/ntrip.service

[Unit]
Description=Ntrip client service


[Service]
SyslogIdentifier=ntrip
ExecStart=bash -c "/usr/local/bin/start_ntrip"
User=pi
KillSignal=SIGINT
Restart=always
RestartSec=15
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
